# scripts/systemd/news-scraper-open.service
[Unit]
Description=News Scraper Open pipeline
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=news-scraper
WorkingDirectory=/opt/news-scraper-open
Environment=CHROME_BIN=/snap/bin/chromium
Environment=CHROMEDRIVER_PATH=/usr/local/bin/chromedriver
EnvironmentFile=/opt/news-scraper-open/.env
ExecStart=/opt/news-scraper-open/.venv/bin/python /opt/news-scraper-open/main.py run
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
